# docker/cron/crontab
# m h  dom mon dow command
0 3 * * * /var/www/current/bin/console proxy:certificate:cron-job